# debian/linapse-configurator.service
[Unit]
Description=Linapse configurator (static web UI)
After=linapse-service.service

[Service]
ExecStart=/usr/bin/python3 -m http.server 7890 --bind 127.0.0.1 --directory /usr/share/linapse/configurator
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
